# letsencrypt@.timer
[Unit]
Description=Renew LetsEncrypt-certificate daily

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
